# deploy/grimoire.service
# systemd USER service for the Grimoire dashboard (API + built frontend) on :8731.
# Install:
#   mkdir -p ~/.config/systemd/user && cp deploy/grimoire.service ~/.config/systemd/user/
#   systemctl --user daemon-reload
#   systemctl --user enable --now grimoire
#   sudo loginctl enable-linger $USER   # one-time: run without an open WSL session
# Logs: journalctl --user -u grimoire -f
[Unit]
Description=Grimoire dashboard (API serves the built frontend)
After=network-online.target

[Service]
WorkingDirectory=/home/adyu/apps/grimoire
ExecStart=/home/adyu/apps/grimoire/.venv/bin/uvicorn grimoire.api:app --host 0.0.0.0 --port 8731
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
